# systemd/talking-flower.service
[Unit]
Description=Talking Flower Voice Assistant
After=picoclaw-gateway.service sound.target
Wants=picoclaw-gateway.service

[Service]
Type=simple
User=plue
WorkingDirectory=/home/plue/talking-flower/voice-assistant
ExecStart=/home/plue/talking-flower/voice-assistant/.venv/bin/python3 voice_assistant.py
Environment=INPUT_MODE=gpio
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
